# ops/divvun-update.service
[Unit]
Description=Nightly Divvun Update

[Service]
Type=oneshot
ExecStart=/home/kevin/setup-linux.sh
